# docker/pass3/to_copy/etc/service/pass3/run
#!/bin/sh

HOME=/releases exec /releases/bin/pass3 foreground