# docker_images/base/derived/ssk_web/runit/ssk_web/run
#!/bin/sh

HOME=/app exec /app/ssk_web/ssk_web_release/bin/ssk_web_release foreground